# userpatches/overlay/services/led_blink.service
[Unit]
Description=LED Blink Service
After=network.target

[Service]
ExecStart=/usr/local/bin/led_blink.sh
Restart=always

[Install]
WantedBy=multi-user.target
